# mirror-console/systemd/mirror-console-web.service
[Unit]
Description=Smart Mirror camera console — web front-end (Express + React)
After=network.target mirror-console-backend.service
Wants=mirror-console-backend.service

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/smartMirror/mirror-console/server
# Full path to node. This is nvm-managed and version-specific — update it if
# the node version changes (`which node` on the Pi), or symlink a stable path.
ExecStart=/home/admin/.nvm/versions/node/v25.9.0/bin/node /home/admin/smartMirror/mirror-console/server/index.js
Environment=PORT=8000
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
