# camera/face_reco.service
[Unit]
Description=Face Recognition Daemon for MMM-Profile
Documentation=https://github.com/Domes711/smartMirror
After=network.target mosquitto.service
Wants=mosquitto.service

[Service]
Type=simple
User=admin
Group=admin
WorkingDirectory=/home/admin/smartMirror/camera
ExecStart=/usr/bin/python3 /home/admin/smartMirror/camera/face_reco_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables (optional overrides)
Environment="MQTT_BROKER=127.0.0.1"
Environment="MQTT_PORT=1883"

[Install]
WantedBy=multi-user.target
